# tools/proxy_watcher/systemd/proxy_watcher.service
[Unit]
Description=bash-kit proxy watcher service
After=network.target
Wants=network.target

[Service]
Type=simple
Restart=always
RestartSec=10
ExecStart=/usr/local/bin/proxy_watcher
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
